# rpi-clock.service
[Unit]
Description=Raspberry Pi Clock Display
After=network-online.target time-sync.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/bin/rpi-clock/clock.py
WorkingDirectory=/usr/bin/rpi-clock
Restart=always
RestartSec=10
User=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target